# mail-service/service-conf/conf/supervisor.conf
[program:mail-service]
command=/usr/bin/java  -Xms2048m -Xmx2048m -XX:PermSize=256m -XX:MaxPermSize=512m  -jar /srv/jar/mail-service/current/mail-provider-1.0-SNAPSHOT.jar
user=tomcat
autostart=false
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/srv/jar/mail-service/current/log/mail-stdout.log
stderr_logfile=/srv/jar/mail-service/current/log/mail-stderr.log
